# build/installer/debpackage/postinst
#!/bin/sh
# postinst script for qnapzfs
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package
OCF_MODULE="/lib/modules/4.2.8/ocf.ko"
CRYPTO_MODULE="/lib/modules/4.2.8/cryptosoft.ko"
CRYPTO_DEV_MODULE="/lib/modules/4.2.8/cryptodev.ko"
ZFS_MODULE="/lib/modules/4.2.8/zfs.ko"
LPL_MODULE="/lib/modules/4.2.8/lpl.ko"
SCST_MODULE="/lib/modules/4.2.8/scst.ko"
ISCSI_SCST_MODULE="/lib/modules/4.2.8/iscsi_scst.ko"
SCST_VDISK_MODULE="/lib/modules/4.2.8/zscst_vdisk.ko"
ISERT_SCST_MODULE="/lib/modules/4.2.8/isert_scst.ko"


ZFS_BINARY="/sbin/zfs"
ZPOOL_BINARY="/sbin/zpool"
ISCSIADM_BINARY="/sbin/iscsiadm"
SCSTADM_BINARY="/sbin/scstadm"
SCSIDEVADM_BINARY="/sbin/scsidevadm"
DAEMON="/sbin/iscsi-scstd"

case "$1" in
    configure)
	chrpath -r '/usr/lib' $ZFS_BINARY
	chrpath -r '/usr/lib' $ZPOOL_BINARY
	echo "---------------------------------------- "
	echo "----------Loading Modules----------"
	insmod $OCF_MODULE
	insmod $CRYPTO_MODULE
	insmod $LPL_MODULE
	insmod $ZFS_MODULE
	insmod $CRYPTO_DEV_MODULE
	#------------------- Load SCST Modules -------------------#
	insmod $SCST_MODULE
	insmod $ISCSI_SCST_MODULE
	insmod $SCST_VDISK_MODULE
	insmod $ISERT_SCST_MODULE
	cp /lib/modules/4.2.8/kernel/drivers/infiniband/core/rdma_cm.ko /lib/modules/4.2.8/
	depmod -a
	echo "----------Successfully loaded all modules----------"
	/sbin/iscsi-scstd &
	echo "---------------------------------------- "
	echo "libcrc32c\nib_cm\niw_cm\nrdma_cm\nocf\ncryptosoft\nlpl\nzfs\ncryptodev\nscst\niscsi_scst\nzscst_vdisk\nisert_scst" >> /etc/modules
	update-rc.d iscsi-scstd defaults

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
